# kas/meta-tezsign/recipes-core/tezsign-core/files/tezsign.service
[Unit]
Description=Runs tezsign
After=setup-gadget.service
Requires=setup-gadget.service

[Service]
Type=simple
User=tezsign
Group=tezsign
Environment="DATA_STORE=/data/tezsign"
Environment="LOG_LEVEL=warn"
ExecStart=/app/tezsign
Restart=on-failure
RestartSec=2

AllowedCPUs=1-3

StandardOutput=null
StandardError=null

ControlGroup=yes
Delegate=yes

[Install]
WantedBy=multi-user.target